# PISailbox/pisailbox.service
[Unit]
Description=PiSailBox GPS Tracker & Camera
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/pisailbox
ExecStart=/usr/bin/python3 /opt/pisailbox/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# UART + SIM7600 need time after power-on (Python side adds extra settle; see PISAILBOX_MODEM_SETTLE)
ExecStartPre=/bin/sleep 15

# Environment — edit SERVER URL here or in /opt/pisailbox/config.py
Environment=PISAILBOX_SERVER=https://lovesailing.ai

[Install]
WantedBy=multi-user.target
